# archlinux-etcd/config/etcd.service
[Unit]
Description=Clustered etcd

[Service]
Restart=always
ExecStart=/usr/bin/etcd

[Install]
WantedBy=multi-user.target
